# bin/tl-capture-once.service
[Unit]
Description=Timelapse: one capture

[Service]
Type=oneshot
ExecStart=/home/sha/timelapse/bin/capture-once.sh
User=sha
Group=sha


