# bin/tl-daily-process.service
[Unit]
Description=Timelapse: daily video build & hour filtering

[Service]
Type=oneshot
ExecStart=/home/sha/timelapse/bin/daily-process.sh
User=pi
Group=pi

